# supervisord.conf
[supervisord]
nodaemon=true

[program:twgts]
command=/usr/local/bin/twgts
directory=/data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
